# sensu-go/files/sensu-agent.in
#!/bin/sh
# $FreeBSD: head/sysutils/sensu-go/files/sensu-agent.in jp $

# PROVIDE: sensu-agent
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable sensu-agent:
#
# sensu_agent_enable (bool):	Set to YES to enable sensu-agent
# 				Default: NO
# sensu_agent_conf (str):	sensu-agent configuration file
#				Default: ${PREFIX}/etc/sensu-agent.yml

. /etc/rc.subr

name=sensu_agent
rcvar=sensu_agent_enable
load_rc_config ${name}

sensu_agent_bin=/usr/local/bin/sensu-agent

: ${sensu_agent_enable:="NO"}
: ${sensu_agent_conf:="-c %%PREFIX%%/etc/sensu-agent.yml"}

command=/usr/sbin/daemon
pidfile="/var/run/sensu-agent.pid"
command_args="-rfP ${pidfile} ${sensu_agent_bin} start ${sensu_agent_conf} >> /var/log/sensu-agent.log 2>&1"

sensu_agent_prestart() {
# Have to empty rc_flags so they don't get passed to daemon(8)
	rc_flags=""
}

run_rc_command "$1"
